# persistence/cron_job_creation.yml
---
title: suspicious_cron_job_creation
description: Detects creation of cron jobs with suspicious commands or by unusual processes.
author: nanos.sh
severity: medium
mode: live
detection_mode: scheduled
schedule: "*/5 * * * *"
lookback: 15m
folder: endpoint
mitre_tactics:
  - TA0003
  - TA0002
mitre_techniques:
  - T1053.003
tags:
  - cron
  - linux
  - persistence
ai_triage_hints:
  ignore_when:
    - "Cron job is from package manager or system update"
    - "Known automation tool creating scheduled task"
  suspicious_when:
    - "Cron command includes network activity (curl, wget)"
    - "Cron command runs encoded or obfuscated script"
---
source_type=linux_audit OR source_type=endpoint_edr
| where action="process_create"
| where process_name="crontab" OR file_path CONTAINS "/cron.d/" OR file_path CONTAINS "/crontab"
| risk score=45 entity=src_host factor="Cron job modified"
| where process CONTAINS "curl" OR process CONTAINS "wget" OR process CONTAINS "nc " OR process CONTAINS "python"
| risk score=75 factor="Network command in cron"
| where process CONTAINS "base64" OR process CONTAINS "|sh" OR process CONTAINS "|bash"
| risk score=85 factor="Encoded or piped execution"
| table timestamp, src_host, user, process, risk_score, risk_factors
